# spacewalk/web/supervisord.conf
[supervisord]
nodaemon=true

[program:spacewalk-setup]
command=/root/config_spacewalk.sh
autorestart=false