# tutoriais/web/supervisor-httpd.ini
[program:httpd]
priority = 100
command = /bin/bash -c "/bin/sleep 2 && /usr/sbin/httpd -c \"ErrorLog /dev/stdout\" -DFOREGROUND"
redirect_stderr = true
stdout_logfile = /var/log/httpd/error_log
stdout_events_enabled = true
startsecs = 5
autorestart = true
